# vendor/dell/s6010/lib/systemd/system/opx-platform-init.service
[Unit]
Description=Platform initialization service
After=dev-i2c-0.device dev-i2c-1.device
Wants=dev-i2c-0.device dev-i2c-1.device
Before=opx-pas.service
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/bin/opx_platform_init.sh
RemainAfterExit=yes

# Resource Limitations
LimitNOFILE=16384

[Install]
WantedBy=multi-user.target
