# vendor/dell/s6010/lib/systemd/system/opx-platform-init.service
[Unit]
Description=Platform initialization service
After=sys-devices-pci0000:00-0000:00:13.0-i2c\x2d0-i2c\x2ddev-i2c\x2d0.device sys-devices-pci0000:00-0000:00:1f.3-i2c\x2d1-i2c\x2ddev-i2c\x2d1.device
Wants=sys-devices-pci0000:00-0000:00:13.0-i2c\x2d0-i2c\x2ddev-i2c\x2d0.device sys-devices-pci0000:00-0000:00:1f.3-i2c\x2d1-i2c\x2ddev-i2c\x2d1.device
Before=opx-pas.service
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/bin/opx_platform_init.sh
RemainAfterExit=yes

# Resource Limitations
LimitNOFILE=16384

[Install]
WantedBy=multi-user.target
