# vendor/dell/z9264f/lib/systemd/system/opx-platform-init.service
[Unit]
Description=Platform initialization service
Before=opx-pas.service
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/bin/opx_platform_init.sh
RemainAfterExit=yes

# Resource Limitations
LimitNOFILE=16384

[Install]
WantedBy=multi-user.target
