# vendor/dell/z9264f/lib/systemd/system/opx-platform-init.service
[Unit]
Description=Platform initialization service
After=sys-devices-i2c\x2d603-i2c\x2ddev-i2c\x2d603.device sys-devices-i2c\x2d604-i2c\x2ddev-i2c\x2d604.device sys-devices-i2c\x2d605-i2c\x2ddev-i2c\x2d605.device sys-devices-i2c\x2d606-i2c\x2ddev-i2c\x2d606.device sys-devices-i2c\x2d607-i2c\x2ddev-i2c\x2d607.device sys-devices-i2c\x2d608-i2c\x2ddev-i2c\x2d608.device sys-devices-i2c\x2d609-i2c\x2ddev-i2c\x2d609.device sys-devices-i2c\x2d610-i2c\x2ddev-i2c\x2d610.device sys-devices-i2c\x2d611-i2c\x2ddev-i2c\x2d611.device
Wants=sys-devices-i2c\x2d603-i2c\x2ddev-i2c\x2d603.device sys-devices-i2c\x2d604-i2c\x2ddev-i2c\x2d604.device sys-devices-i2c\x2d605-i2c\x2ddev-i2c\x2d605.device sys-devices-i2c\x2d606-i2c\x2ddev-i2c\x2d606.device sys-devices-i2c\x2d607-i2c\x2ddev-i2c\x2d607.device sys-devices-i2c\x2d608-i2c\x2ddev-i2c\x2d608.device sys-devices-i2c\x2d609-i2c\x2ddev-i2c\x2d609.device sys-devices-i2c\x2d610-i2c\x2ddev-i2c\x2d610.device sys-devices-i2c\x2d611-i2c\x2ddev-i2c\x2d611.device
Before=opx-pas.service
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/bin/opx_platform_init.sh
RemainAfterExit=yes

# Resource Limitations
LimitNOFILE=16384

[Install]
WantedBy=multi-user.target
